# system/zyl-credential.service
[Unit]
Description=Zyl OS Credential Manager
After=dbus.service zyl-compositor.service
Wants=zyl-compositor.service

[Service]
Type=dbus
BusName=org.zylos.CredentialManager
ExecStart=/usr/bin/zyl-credential
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening — handles encryption keys and credentials
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
