# system/zyl-alarm.service
[Unit]
Description=Zyl OS Alarm/Scheduler Service
After=dbus.service

[Service]
Type=dbus
BusName=org.zylos.AlarmService
ExecStart=/usr/bin/zyl-alarm
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
