# system/zyl-accessibility.service
[Unit]
Description=Zyl OS Accessibility Service
After=dbus.service

[Service]
Type=dbus
BusName=org.zylos.Accessibility
ExecStart=/usr/bin/zyl-accessibility
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
