# system/zyl-power.service
[Unit]
Description=Zyl OS Power Management Service
After=dbus.service zyl-compositor.service
Wants=zyl-compositor.service

[Service]
Type=dbus
BusName=org.zylos.PowerManager
ExecStart=/usr/bin/zyl-power
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
