# system/zyl-sensors.service
[Unit]
Description=Zyl OS Sensor Service
After=dbus.service

[Service]
Type=dbus
BusName=org.zylos.SensorService
ExecStart=/usr/bin/zyl-sensors
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
