# system/zyl-bluetooth.service
[Unit]
Description=Zyl OS Bluetooth Service (BlueZ)
After=bluetooth.service

[Service]
Type=dbus
BusName=org.zylos.BluetoothService
ExecStart=/usr/bin/zyl-bluetooth
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=zyl-os.target
