# system/zyl-input.service
[Unit]
Description=Zyl OS Input Service (keyboard, touch, hardware keys)
After=dbus.service
Requires=dbus.service
PartOf=zyl-os.target

[Service]
Type=dbus
BusName=org.zylos.InputService
ExecStart=/usr/bin/zyl-input
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=zyl-os.target
