# system/zyl-wam.service
[Unit]
Description=Zyl OS Web Application Manager
After=zyl-compositor.service
Requires=zyl-compositor.service
PartOf=zyl-compositor.service

[Service]
Type=dbus
BusName=org.zylos.WebAppManager
User=root
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/0
ExecStart=/usr/bin/zyl-wam
Restart=on-failure
RestartSec=2

# Security hardening — root required for Wayland/XDG_RUNTIME_DIR session
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
# ProtectControlGroups=false — WAM OOM killer가 cgroup v2로 앱별 메모리 제한 설정
ProtectControlGroups=false
RestrictNamespaces=true
# OOM killer cgroup 경로 접근 허용
ReadWritePaths=/sys/fs/cgroup/zylos.apps

[Install]
WantedBy=graphical.target
