# system/zyl-compositor.service
[Unit]
Description=Zyl OS Wayland Compositor
After=systemd-logind.service plymouth-quit.service
Wants=systemd-logind.service
ConditionPathExists=/dev/dri/card0

[Service]
Type=simple
User=root
Environment=WLR_BACKENDS=drm
Environment=WLR_DRM_DEVICES=/dev/dri/card0
Environment=XDG_RUNTIME_DIR=/run/user/0
ExecStart=/usr/bin/zyl-compositor
Restart=on-failure
RestartSec=3

# Security hardening — root required for DRM/logind session
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true

[Install]
WantedBy=graphical.target
