# system/zyl-usb.service
[Unit]
Description=Zyl OS USB/MTP Manager Service
After=dbus.service

[Service]
Type=dbus
BusName=org.zylos.UsbManager
ExecStart=/usr/bin/zyl-usb
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
