# system/zyl-notification.service
[Unit]
Description=Zyl OS Notification Service
After=dbus.service

[Service]
Type=dbus
BusName=org.zylos.Notification
ExecStart=/usr/bin/zyl-notification
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
