# system/zyl-wam.service
[Unit]
Description=Zyl OS Web Application Manager
After=zyl-compositor.service
Requires=zyl-compositor.service
PartOf=zyl-compositor.service

[Service]
Type=dbus
BusName=org.zylos.WebAppManager

# Run as dedicated non-root service account instead of root.
# The zyl-wam user must be created by the installer/package:
#   useradd -r -s /usr/sbin/nologin -d /var/lib/zyl-os/wam -G video,render zyl-wam
# XDG_RUNTIME_DIR is set by pam_systemd or the RuntimeDirectory directive.
User=zyl-wam
Group=zyl-wam
# video/render groups required for DRI/GPU access
SupplementaryGroups=video render

# Wayland display socket access: zyl-compositor must have granted group access
# (or use /run/zyl-wam/wayland-0 scoped to zyl-wam user)
Environment=WAYLAND_DISPLAY=wayland-0
RuntimeDirectory=zyl-wam
RuntimeDirectoryMode=0700
# XDG_RUNTIME_DIR must match the RuntimeDirectory path so WebKit/Wayland
# socket lookup works correctly under the zyl-wam service account.
Environment=XDG_RUNTIME_DIR=/run/zyl-wam

ExecStart=/usr/bin/zyl-wam
Restart=on-failure
RestartSec=2

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
# WAM OOM killer requires cgroup v2 write access
ProtectControlGroups=false
RestrictNamespaces=true
# Allow cgroup path for OOM killer app memory accounting
ReadWritePaths=/sys/fs/cgroup/zylos.apps
# App data directories
ReadWritePaths=/data/apps /var/lib/zyl-os/wam /var/cache/zyl-os/wam

[Install]
WantedBy=graphical.target
