# system/zyl-camera.service
[Unit]
Description=Zyl OS Camera Service (V4L2)
After=zyl-wam.service
Requires=zyl-wam.service

[Service]
Type=dbus
BusName=org.zylos.CameraService
ExecStart=/usr/bin/zyl-camera
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/data/camera
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=zyl-os.target
