# system/zyl-apparmor.service
[Unit]
Description=Zyl OS AppArmor Profile Loader
DefaultDependencies=no
Before=zyl-wam.service zyl-compositor.service
After=apparmor.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Fail closed on parser errors instead of silently ignoring broken profiles.
# Missing glob matches are ignored safely; individual profile load failures abort.
ExecStart=/bin/sh -eu -c '
  found=0; \
  for p in /etc/apparmor.d/zyl-*; do \
    [ -e "$p" ] || continue; \
    found=1; \
    apparmor_parser -r "$p"; \
  done; \
  if [ "$found" -eq 0 ]; then \
    echo "No ZylOS AppArmor profiles found under /etc/apparmor.d" >&2; \
  fi'

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=zyl-os.target
