# system/zyl-logger.service
[Unit]
Description=Zyl OS Logger Service
After=dbus.service local-fs.target

[Service]
Type=dbus
BusName=org.zylos.Logger
ExecStart=/usr/bin/zyl-logger
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/zyl-os
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
