# system/zyl-wifi.service
[Unit]
Description=Zyl OS WiFi Service (NetworkManager)
After=NetworkManager.service

[Service]
Type=dbus
BusName=org.zylos.WifiService
ExecStart=/usr/bin/zyl-wifi
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/data/vpn
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=zyl-os.target
